# bots/investBot/com.investbot.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.investbot.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/npx</string>
        <string>tsx</string>
        <string>/Users/jasonlee/codingProjects/investBot/src/weeklyReport.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jasonlee/codingProjects/investBot</string>

    <!-- Run every Sunday at 6:00 PM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jasonlee/codingProjects/investBot/logs/launchd-weekly.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jasonlee/codingProjects/investBot/logs/launchd-weekly-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
