# bots/investBot/com.investbot.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.investbot.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/npx</string>
        <string>tsx</string>
        <string>/Users/jasonlee/codingProjects/investBot/src/dailyRunner.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jasonlee/codingProjects/investBot</string>

    <!-- Run at 9:35 AM ET (market opens 9:30, wait 5 min for stability) -->
    <!-- Adjust hour based on your timezone -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
        <!-- Tuesday -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
        <!-- Wednesday -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
        <!-- Thursday -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
        <!-- Friday -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>35</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jasonlee/codingProjects/investBot/logs/launchd-daily.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jasonlee/codingProjects/investBot/logs/launchd-daily-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Keep trying if it fails -->
    <key>KeepAlive</key>
    <false/>

    <!-- Run on login/wake to catch up if 9:35 was missed (e.g. laptop was asleep) -->
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
